# 02-builtin-server/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python /var/www/app/application.py
